# refresh-hugo-amritanshu-in.timer
[Unit]
Description=Refresh hugo amritanshu in

[Timer]
OnBootSec=5min
OnUnitActiveSec=30s

[Install]
WantedBy=timers.target